# resources/etc/supervisor/conf.d/supervisord_nepi.conf
[supervisord]
nodaemon=false

[program:nepi_engine]
command=/bin/bash /opt/nepi/scripts/nepi_engine_start
autostart=true
autorestart=false

[program:nepi_rui]
command=/bin/bash /opt/nepi/scripts/nepi_rui_start
autostart=true
autorestart=false

[program:nepi_license]
command=/bin/bash /opt/nepi/scripts/opt/nepi/etc/license/nepi_check_license_start
autostart=true
autorestart=true

[program:nepi_storage]
command=/bin/bash /opt/nepi/scripts/nepi_storage_start
autostart=true
autorestart=true

[program:nepi_sync]
command=/bin/bash /opt/nepi/scripts/nepi_sync_start
autostart=true
autorestart=false

